# roles/base/templates/etc/profile.d/pvc.sh.j2
#!/bin/sh

# Message of the day script to print PVC cluster status
# {{ ansible_managed }}

echo
pvc -q cluster status -f short
pvc -q node list
