# roles/base/templates/etc/profile.d/w.sh.j2
#!/bin/sh

# Message of the day script to print active users
# {{ ansible_managed }}

export PROCPS_FROMLEN=36 PROCPS_USERLEN=12
w
