# roles/pvc/templates/cputuning/user.slice.j2
# PVC VM slice unit
# {{ ansible_managed }}

[Unit]
Description = User and Session Slice
Documentation = man:systemd.special(7)
Before = slices.target

[Slice]
CPUAccounting = true
Delegate = true
CPUAffinity = {{ cpuset_system }}
AllowedCPUs = {{ cpuset_system }}
