# roles/pvc/templates/automirror/pvc-automirror.timer.j2
[Unit]
Description=[Timer] PVC VM automirror

[Timer]
Unit=pvc-automirror.service
OnCalendar={{ pvc_automirror.schedule.time }}

[Install]
WantedBy=pvc.target
