# roles/pvc/templates/autobackup/pvc-autobackup-full.timer.j2
[Unit]
Description=[Timer] PVC VM autobackup (forced-full)

[Timer]
Unit=pvc-autobackup-full.service
OnCalendar={{ pvc_autobackup.schedule.forced_full_time }}

[Install]
WantedBy=pvc.target
