# roles/pvc/templates/autobackup/pvc-autobackup-normal.service.j2
[Unit]
Description=[Cron] PVC VM autobackup (normal)

[Service]
Type=oneshot
IgnoreSIGPIPE=false
KillMode=process
ExecStart=/usr/bin/pvc --quiet vm autobackup --cron
