# roles/base/templates/etc/cron.d/update-motd.j2
# cron file for motd
# {{ ansible_managed }}

@reboot  root  /bin/sh /usr/local/sbin/update-motd.sh &>/dev/null
