# roles/base/templates/etc/cron.d/update-motd.j2
# cron file for motd
# {{ ansible_managed }}

PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
@reboot  root  /usr/local/sbin/update-motd.sh &>/dev/null
