# roles/base/templates/etc/cron.d/update-issue.j2
# cron file for update-issue
# {{ ansible_managed }}

PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
@reboot  root  /usr/local/sbin/update-issue.sh &>/dev/null
